# deb/prerm
#! /bin/bash

systemctl stop a4wd2-roscore.service
systemctl disable a4wd2-roscore.service
